# firefox-frontend.service
[Unit]
Description=HAPi Dashboard Frontend - Firefox Display
After=network-online.target
Wants=network-online.target
Documentation=https://github.com/yourusername/HAPi-Dashboard

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/HAPi-Dashboard
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"
Environment="PYTHONUNBUFFERED=1"
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/python3 -u /home/pi/HAPi-Dashboard/frontend-firefox.py --url http://YOUR-SERVER-IP:5000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=HAPi-Dashboard-Frontend
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
