# HAPi-Dashboard.service
[Unit]
Description=HAPi Dashboard - Home Assistant Family Dashboard
After=network-online.target display-manager.service
Wants=network-online.target
Documentation=https://github.com/yourusername/HAPi-Dashboard

[Service]
Type=simple
User=root
WorkingDirectory=/home/dietpi/HAPi-Dashboard
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/root/.Xauthority"
Environment="PYTHONUNBUFFERED=1"
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/python3 -u /home/dietpi/HAPi-Dashboard/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=HAPi-Dashboard
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target